# packages/addons/service/pijuice/system.d/pijuice.service
[Unit]
Description=PiJuice status service
Requires=dev-i2c\x2d1.device
After=dev-i2c\x2d1.device

[Service]
Type=idle
# was pijuice
User=kodi
WorkingDirectory=/var/lib/pijuice/
ExecStart=/storage/.kodi/addons/service.pijuice/bin/pijuice_sys.py
ExecStopPost=/usr/bin/pijuice_sys.py stop
Restart=always

[Install]
WantedBy=multi-user.target

ExecStart=usbmuxd --systemd
